# vrc-join-notify.service
[Unit]
Description=VRChat Join Notifier (user)
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/vrc_join_notify.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target

